# conf/supervisord/button.conf
[program:botton]
directory=/home/botton/botton
command=/usr/bin/python -u main.py
stderr_logfile=/var/log/botton.log
stdout_logfile=/var/log/botton.log
autostart=true
autorestart=true
